# deploy/manager-restart.service
[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart manager

[Install]
WantedBy=multi-user.target
